# ops/launchagents/com.mrtyndall.home-base-mcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mrtyndall.home-base-mcp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/npm</string>
    <string>run</string>
    <string>mcp:http</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/matt/projects/home-base</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>MCP_PORT</key>
    <string>8081</string>
    <key>HOME_BASE_API_URL</key>
    <string>http://127.0.0.1:3002/api/v1</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/home-base-mcp.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/home-base-mcp.log</string>
</dict>
</plist>
